# omen-fan/systemd/omen-fan.service
[Unit]
Description=Omen Fan Control Service
After=multi-user.target

[Service]
ExecStart=/path/to/your/project/target/release/omen-fan
Restart=always
User=root

[Install]
WantedBy=multi-user.target